# repo/ops/nightshift.service
[Unit]
Description=Claude Night Shift runner
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/bin/bash __HOME__/claude-night-shift/night-shift.sh
